# masked-emails-daemon/DEBIAN/postinst
#!/bin/sh

chmod 0755 /etc/init.d/masked-emails

systemctl daemon-reload
/etc/init.d/masked-emails start
update-rc.d masked-emails defaults
