# masked-emails-daemon/DEBIAN/prerm
#!/bin/sh
#!/bin/sh

/etc/init.d/masked-emails stop
update-rc.d -f masked-emails remove
